# package/rpm/usr/lib/systemd/system/udpac.service
[Unit]
Description=udpack
After=network.target

[Service]
ExecStart=/usr/sbin/udpac
Type=simple

[Install]
WantedBy=default.target
